# astroquant_orchestrator.service
[Unit]
Description=AstroQuant Signal Orchestrator
After=network.target astroquant_tradingbot.service chrome_remote_debug.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/workspaces/newcpu
ExecStart=/workspaces/newcpu/.venv/bin/python3 /workspaces/newcpu/start_astroquant.py
Restart=always
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target